# templates/init/Debian.init.erb
#! /bin/sh
### BEGIN INIT INFO
# Provides:          teamspeak
# Required-Start:    networking
# Required-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:      S 0 1 6
# Short-Description: TeamSpeak Server Daemon
# Description:       Starts/Stops/Restarts the TeamSpeak Server Daemon
### END INIT INFO

set -e

PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
user=teamspeak
dir=/opt/teamspeak
exec=./ts3server_startscript.sh

test -x $DAEMON || exit 0

case "$1" in
  start)
    su - $user -c "cd $dir;$exec start"
  ;;
  stop)
    su - $user -c "cd $dir;$exec stop"
  ;;
  restart)
    su - $user -c "cd $dir;$exec restart"
  ;;
  status)
    su - $user -c "cd $dir;$exec status"
  ;;
  *)
    echo $"Usage: $0 {start|stop|restart|status}"
    exit 2
esac
exit 0 